# packages/appserver-service/carbonio-mailbox-internal-api-sidecar.service
[Unit]
Description=Carbonio Mailbox Internal Api Sidecar
Documentation=https://docs.zextras.com/
Requires=network-online.target
After=network-online.target service-discover.service carbonio-appserver.service
PartOf=service-discover.target

[Service]
User=carbonio-mailbox
Type=notify
NotifyAccess=all
ExecStart=/usr/bin/service-discover-wrapper.sh \
    -token-file /etc/carbonio/mailbox-internal-api/service-discover/token \
    -admin-bind localhost:0 \
    -sidecar-for carbonio-mailbox-internal-api
ExecReload=/usr/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=15s
KillSignal=SIGINT
LimitNOFILE=65536

# Hardening
PrivateTmp=yes
ProtectSystem=strict
NoNewPrivileges=yes
PrivateDevices=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes
ProtectHostname=yes
ProtectClock=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=service-discover.target
