# packages/appserver-db/carbonio-appserver-db.service
[Unit]
Description=Carbonio Appserver DB Daemon
PartOf=carbonio-appserver.target

[Service]
User=zextras
ExecStart=/opt/zextras/common/bin/mysqld_safe \
  --defaults-file=/opt/zextras/conf/my.cnf \
  --external-locking \
  --log-error=/opt/zextras/log/mysql_error.log \
  --malloc-lib=/opt/zextras/common/lib/libjemalloc.so \
  --ledir=/opt/zextras/common/sbin
LimitNOFILE=524288
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=carbonio-appserver.target
