# packages/directory-server/carbonio-directory-server.target
[Unit]
Description=Carbonio Directory Server
Documentation=man:systemd.special(7)
Requires=multi-user.target
After=multi-user.target
PartOf=carbonio-ce.target carbonio.target
AllowIsolate=yes

[Install]
WantedBy=multi-user.target